# phlow-nfdump/phlow-nfdump-src/src/etc/cron.hourly/nfcapd.cron
#!/bin/bash

/usr/bin/nfexpire -e /var/log/nfcapd -t 6d
