# phlow-nfdump/phlow-nfdump-src/src/etc/init.d/nfcapd
#!/bin/sh
# 
### BEGIN INIT INFO
# Provides:          nfcapd
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: netflow capture daemon
# Description: 	     netflow capture daemon
# chkconfig: 345 99 01
# processname: nfcapd
# File: nfcapd
### END INIT INFO

#
# Author : Michael Fung (http://3open.org/)
# 

# Source function library.
. /etc/init.d/functions

# Source an auxiliary options file if we have one, and pick up DAEMON_OPTIONS.
[ -r /etc/sysconfig/nfcapd ] && . /etc/sysconfig/nfcapd

DAEMON=nfcapd
LOCK_FILE="/var/lock/subsys/$DAEMON"
RETVAL=0

DAEMON1=nfdump-ascii
LOCK_FILE1="/var/lock/subsys/$DAEMON1"
RETVAL1=0
  
status_nfcapd () {
	if pidof $DAEMON > /dev/null 2>&1; then
		return 0
	else
	  return 1
	fi
}

status_nfdump-ascii() {
	if pidof $DAEMON1 > /dev/null 2>&1; then
                return 0
        else
          return 1
        fi
}

do_start() {
	status_nfcapd
	if [ $? -eq 0 ]; then
		echo "$DAEMON already started..."
		return 1
	fi
	
	echo -n $"Starting $DAEMON: "
	daemon /usr/bin/nfcapd -D $DAEMON_OPTIONS
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && touch $LOCK_FILE
	return $RETVAL

	status_nfdump-ascii
        if [ $? -eq 0 ]; then
                echo "$DAEMON1 already started..."
                return 1
        fi

        echo -n $"Starting $DAEMON1: "
        daemon /usr/bin/nfdump-ascii.sh
        RETVAL1=$?
        echo
        [ $RETVAL1 -eq 0 ] && touch $LOCK_FILE1
        return $RETVAL1
}

do_stop() {
    echo -n $"Stopping $DAEMON: "
    killproc $DAEMON
    RETVAL=$?
    [ $RETVAL -eq 0 ] && rm -f $LOCK_FILE
    echo
    return $RETVAL

    echo -n $"Stopping $DAEMON1: "
    killproc $DAEMON1
    RETVAL1=$?
    [ $RETVAL1 -eq 0 ] && rm -f $LOCK_FILE1
    echo
    return $RETVAL1
}
                
# See how we were called.
case "$1" in
	start)
    do_start
		;;
	
	stop)
		do_stop
    ;;
    
  status)
		status $DAEMON
    ;;

	restart)
		do_stop
		do_start
		;;

	condrestart|try-restart)
                status_nfcapd
                if [ $? -eq 0 ]; then
                        do_stop
                        do_start
                else
                        echo "daemon not running..."
                fi
                ;;

	*)
		echo "Usage: $0 {start|stop|restart|status|condrestart|try-restart}"
		exit 2
		;;

esac
  
exit $RETVAL
  
# End of this script
